# data/supervisor/openttd.conf
[program:openttd]
process_name=%(program_name)s
command = /usr/share/games/openttd/openttd -D
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/openttd/stdout.log

